# infra/launchd/com.quant.daily-market-backfill.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quant.daily-market-backfill</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__REPO_ROOT__/backend/utils/run_daily_market_backfill.py</string>
    <string>--lookback-days</string>
    <string>14</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/daily-market-backfill.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/daily-market-backfill.err.log</string>
</dict>
</plist>
